# rootfs/etc/inittab
::sysinit:/sbin/init
tty1::respawn:/sbin/getty 38400 tty1
ttyS0::respawn:/sbin/getty -L 38400 ttyS0 vt100
